# deploy/abetterme-api.service
[Unit]
Description=A Better Me API (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/ABetterMe
Environment=ABETTERME_HOST=127.0.0.1
Environment=ABETTERME_PORT=8787
Environment=ABETTERME_DB=/root/abetterme-data.db
ExecStart=/root/ABetterMe/.venv/bin/python server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
